# scripts/supervisor/conf.d/leaps.conf
[program:leaps]
command = /opt/leaps/bin/leaps
autostart = true
stdout_logfile = /var/log/leaps-out.log
stderr_logfile = /var/log/leaps-err.log
